# container-files/etc/services.d/selenium/run
#!/usr/bin/with-contenv /usr/bin/execlineb
/usr/bin/java -jar /usr/bin/selenium-server-standalone.jar
